# data/systemd/systemd-timesyncd.service.d/disable-timesyncd-with-installer.conf
[Unit]
# Don't run timesyncd if installer is installed
ConditionFileIsExecutable=!/usr/bin/deepin-installer
